# companies/LawGens/legacy-audit/legacy-services/systemd/lawgens-contract-os.service
[Unit]
Description=LawGens Contract OS
After=network.target mongodb.service redis.service

[Service]
Type=simple
User=lawgens
WorkingDirectory=/opt/lawgens/contract-os
Environment=NODE_ENV=production
Environment=PORT=4190
EnvironmentFile=/opt/lawgens/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawgens-contract-os

[Install]
WantedBy=multi-user.target